# bin/cleantemp.cron
#!/opt/bb/bin/bash

if [ "$1" = '--verbose' ]
then
    rmv='-v'
    echo "Output for cleantemp.cron, which cleans /tmp of files and subdirectories older than $(date --date='2 days ago')."
fi

for i in $(find /tmp -depth -mtime +1)
do
    if [[ ! $i =~ systemd ]]
    then
        [ "$1" = '--verbose' ] && ls -la $i
        \rm $rmv -f $i
    fi
done
