# bin/cleantemp.cron
#!/opt/bb/bin/bash

date

echo 'Log for cleantemp.cron, which cleans /tmp of files older than two days.'

. ${HOME}/.bash_profile

for i in $(find /tmp -atime +1)
do
    \rm -v -rf $i
done
