# bin/kerbcheck.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

tdiff=$(toolkit_expires_in 90000) ## 60 * 60 * 25
expired=$?
if((expired))
then
    toolkit-file-check
    echo -n "Kerberos token "
    if ((tdiff <= 0))
    then
        echo -n "has expired. "
    else
        echo -n "will expire soon. "
    fi
    echo "Run 'toolkit-file-maint'."
fi
