# bin/kerbcheck.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

toolkit-is-expired
expired=$?
if((expired))
then
    echo "Kerberos token has expired. Run 'toolkit-file-maint'."
fi
