# bin/rsyncvm.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

if [ -r /tmp/rsyncvm.cron.pause ]
then
    echo "Sync paused:"
    ls -la /tmp/rsyncvm.cron.pause
else
    rsyncvm "$@"
fi
