# bin/rsyncvm.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

if [ -r /tmp/vmrsync.cron.pause ]
then
    echo "Sync paused:"
    ls -la /tmp/vmrsync.cron.pause
else
    vmrsync "$@"
fi
