# bin/rsyncvm_spacecheck.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

set-nfsmachine; status=$?; [ ! "$status" = 0 ] && exit $status

ssh $NFSMACHINE 'date;/bbsrc/bin/fsusage | grep home; date;echo Total space used;df -k .; date;cd vmstashes;echo vmstashes space used:;du -ks .; date'
