# bin/vmrsync.cron
#!/opt/bb/bin/bash

date

echo 'Log for vmrsync.cron, which rsyncs this vm back to an NFS partition in'
echo 'order to take advantage of NFS snapshoting.'

. ${HOME}/.bash_profile

vmrsync
