# bin/vmrsync.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

vmrsync "$@"
