# bin/vmrsync_spacecheck.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

ssh $NFSMACHINE 'cd vmstashes;echo Total space used:;df -k .;echo vmstashes space used:;du -ks .'
