# bin/vmupdate.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

yum-update -y
apt-update -y
apt-upgrade -y --force-yes
