# bin/vmupdate.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

echo "*** yum ***"
yum-update -y
echo "*** apt ***"
apt-update -y
apt-upgrade -y --force-yes
