# bin/vmupdate.cron
#!/opt/bb/bin/bash

. ${HOME}/.bash_profile

echo "*** yum ***"
yum-update -y
echo status=$?
echo "*** apt ***"
echo "* update"
apt-update -y
echo status=$?
echo "* upgrade"
apt-upgrade -y --force-yes
echo status=$?
